# github520/update.timer
[Unit]
Description=Run GitHub520 hosts update every hour

[Timer]
OnBootSec=5min
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
